# lib/runner.service.tmpl
[Unit]
Description=GitHub Actions self-hosted runner (@RDIR@)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=@RDIR@
ExecStart=@RDIR@/run.sh
Restart=always
RestartSec=5
# default KillMode=control-group: stop kills run.sh AND its Runner.Listener child,
# so no orphaned process keeps holding the GitHub runner session.
TimeoutStopSec=30

[Install]
WantedBy=default.target
