# sites/bitforbytes/units/frontend.service.tmpl
[Unit]
Description=VeriLog frontend (vite preview :@PORT@)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=@HOME@/verilog-src/frontend
Environment=NODE_ENV=production
Environment=PORT=@PORT@
Environment=PATH=@BIN@:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env npm run preview
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
